# ansible/roles/deploy/files/rove.service
[Unit]
Description=rove service

[Service]
User=rove
Group=rove
ExecStart=/usr/local/bin/rove

[Install]
WantedBy=multi-user.target
